# deploy/smartfinances.service
[Unit]
Description=SmartFinances Telegram Bot
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/smartfinances
Environment="PATH=/var/www/smartfinances/venv/bin"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/var/www/smartfinances/google_credentials.json"

# Kill any existing bot processes before starting
ExecStartPre=-/bin/bash -c '/usr/bin/pkill -9 -f "python.*main.py"; exit 0'
ExecStartPre=/bin/sleep 2

ExecStart=/var/www/smartfinances/venv/bin/python3 /var/www/smartfinances/main.py

# Ensure all processes are killed on stop
ExecStopPost=-/bin/bash -c '/usr/bin/pkill -9 -f "python.*main.py"; exit 0'

Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=10

# Logging
StandardOutput=append:/var/www/smartfinances/bot.log
StandardError=append:/var/www/smartfinances/bot.log

[Install]
WantedBy=multi-user.target

